# systemd/pm-frontend.user.service
[Unit]
Description=pm-director Frontend (Vben port 5777)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m http.server 5777 --bind 0.0.0.0
WorkingDirectory=/home/samuel/ruoyi-office-vben/apps/web-antd/dist
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
